# etc/xdg/pardus/autostart/thunar-bookmarks.desktop
[Desktop Entry]
Name=Thunar Bookmarks
Exec=/usr/share/pardus/xfce-settings/thunar-bookmarks
