# etc/xdg/pardus/autostart/xdg-pardus-templates.desktop
[Desktop Entry]
Name=Xfce Templates
Exec=/usr/share/pardus/xfce-settings/xdg-pardus-templates/xdg-pardus-templates
